# deploy/supervisor/photogongju-node.conf
; ╔══════════════════════════════════════════════════════════════════╗
; ║  PhotoGongju Node.js 主站 — Supervisor 进程守护配置               ║
; ║  安装方式:                                                        ║
; ║     sudo cp photogongju-node.conf /etc/supervisor/conf.d/         ║
; ║     sudo supervisorctl reread                                     ║
; ║     sudo supervisorctl update                                     ║
; ║     sudo supervisorctl start photogongju-node                     ║
; ║  管理命令:                                                        ║
; ║     supervisorctl status photogongju-node     查看状态             ║
; ║     supervisorctl restart photogongju-node    重启服务             ║
; ║     supervisorctl tail -f photogongju-node    实时日志             ║
; ╚══════════════════════════════════════════════════════════════════╝

[program:photogongju-node]

; ── 进程标识 ──
process_name=photogongju-node

; ── 启动命令 ──
; 使用 node 直接运行 Express 应用
command=/usr/bin/node /opt/photogongju/node_web/app.js

; 工作目录
directory=/opt/photogongju/node_web

; ── 运行用户 ──
; 注意: 不要用 root，使用 www-data 或新建专用用户更安全
user=www-data

; ── 自动启动 & 崩溃重启 ──
autostart=true              ; 随 supervisor 启动自动运行
autorestart=true            ; 进程意外退出时自动重启
startsecs=5                 ; 启动后在 N 秒内未退出则认为启动成功
startretries=10             ; 启动失败最多重试次数
exitcodes=0,2               ; 哪些退出码视为正常（不触发错误计数）

; ── 优雅停止 ──
stopsignal=TERM             ; 先发送 SIGTERM 给进程
stopwaitsecs=10             ; 等待 10 秒后若未退出则 SIGKILL
stopasgroup=true            ; 同时停止子进程（如子进程池）
killasgroup=true            ; kill 时也发送给进程组

; ── 日志 ──
redirect_stderr=true                                ; stderr 合并到 stdout
stdout_logfile=/opt/photogongju/logs/node_web.log  ; 日志文件路径
stdout_logfile_maxbytes=50MB                        ; 单个日志文件最大 50MB
stdout_logfile_backups=10                           ; 保留最近 10 个日志文件
stdout_capture_maxbytes=0                           ; 不限制 capture 大小

; ── 环境变量 ──
environment=
    NODE_ENV="production",
    PORT="3000",
    AI_SERVICE_URL="http://127.0.0.1:8001"

; ── 系统资源限制 ──
priority=10                 ; 进程优先级（数字越大优先级越低）
