# deploy/supervisor/photogongju-python.conf
; ╔══════════════════════════════════════════════════════════════════╗
; ║  PhotoGongju Python AI 微服务 — Supervisor 进程守护配置          ║
; ║  安装方式:                                                        ║
; ║     sudo cp photogongju-python.conf /etc/supervisor/conf.d/       ║
; ║     sudo supervisorctl reread                                     ║
; ║     sudo supervisorctl update                                     ║
; ║     sudo supervisorctl start photogongju-python                   ║
; ║  管理命令:                                                        ║
; ║     supervisorctl status photogongju-python    查看状态           ║
; ║     supervisorctl restart photogongju-python   重启服务           ║
; ║     supervisorctl tail -f photogongju-python   实时日志           ║
; ╚══════════════════════════════════════════════════════════════════╝

[program:photogongju-python]

; ── 进程标识 ──
process_name=photogongju-python

; ── 启动命令 ──
; 使用虚拟环境中的 uvicorn，生产模式 2 个 worker
; --workers 2: 利用多核 CPU，处理并发请求
; --log-level info: INFO 级别日志
command=/opt/photogongju/python_ai/venv/bin/uvicorn
    main:app
    --host 127.0.0.1
    --port 8001
    --workers 2
    --log-level info
    --timeout-keep-alive 60
    --limit-concurrency 50
    --backlog 128

; 工作目录
directory=/opt/photogongju/python_ai

; ── 运行用户 ──
user=www-data

; ── 自动启动 & 崩溃重启 ──
autostart=true
autorestart=true
startsecs=10                ; AI 服务启动较慢（加载模型等），给更长时间
startretries=10
exitcodes=0,2

; ── 优雅停止 ──
; AI 服务需要更多时间来处理完当前请求并释放 GPU/内存资源
stopsignal=TERM
stopwaitsecs=15             ; 等待 15 秒
stopasgroup=true
killasgroup=true

; ── 日志 ──
redirect_stderr=true
stdout_logfile=/opt/photogongju/logs/python_ai.log
stdout_logfile_maxbytes=50MB
stdout_logfile_backups=10
stdout_capture_maxbytes=0

; ── 环境变量 ──
environment=
    AI_SERVICE_HOST="127.0.0.1",
    AI_SERVICE_PORT="8001",
    OMP_NUM_THREADS="2"         ; 限制 OpenMP 线程数（避免与 worker 争抢 CPU）

; ── 系统资源 ──
priority=20                 ; 比 Node 优先级低（Node 负责用户请求入口）
