# deploy/supervisor/python-cyberphoto.conf
[program:cyberphoto-python]
command=/www/wwwroot/photogongju/python-service/venv/bin/python /www/wwwroot/photogongju/python-service/main.py
directory=/www/wwwroot/photogongju/python-service
autostart=true
autorestart=true
user=www
environment=PYTHONUNBUFFERED=1
stdout_logfile=/www/wwwroot/photogongju/logs/python.log
stderr_logfile=/www/wwwroot/photogongju/logs/python-error.log
